# autostart/.config/autostart/blur-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Wallpaper Blur
Comment=Creates a blurred copy of a background image
Exec=$HOME/.local/bin/blur-wallpaper.sh
X-GNOME-Autostart-enabled=true
NoDisplay=true
Categories=Utility;
